# scripts/launchd/com.tak.minitools.youtube-mail-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tak.minitools.youtube-mail-digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>youtube-mail-digest</string>
        <string>--hours</string>
        <string>24</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/tak/Projects/minitools</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/tak/.asdf/shims:/opt/homebrew/bin</string>
    </dict>
    <!--
      注意: 初回は必ず対話的に `uv run youtube-mail-digest` を一度実行し、
      Gmail OAuth ブラウザ認証を済ませて token.pickle を生成してから launchd 登録すること。
      未認証のまま launchd で起動すると run_local_server がブラウザを開けず永久ハングする。
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/tak/Projects/minitools/logs/youtube-mail-digest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tak/Projects/minitools/logs/youtube-mail-digest.log</string>
</dict>
</plist>
